# share/meta-pluto/meta-sdk/recipes-pluto/pluto-init/files/pluto-init.service
[Unit]
Description=Startup script for Pluto applications.

[Service]
ExecStart=/usr/bin/pluto-init-script

[Install]
WantedBy=multi-user.target